# root/etc/s6-overlay/s6-rc.d/set-lib-cwa-perms/run
#!/bin/bash

echo "Recursively setting ownership of everything in /calibre-library and /app/calibre-web-automated to abc:abc..."
chown -R abc:abc /calibre-library
chown -R abc:abc /app/calibre-web-automated

