# root/etc/s6-overlay/s6-rc.d/calibre-binaries-setup/up
/etc/s6-overlay/s6-rc.d/calibre-binaries-setup/run