# root/etc/s6-overlay/s6-rc.d/universal-calibre-setup/up
/etc/s6-overlay/s6-rc.d/universal-calibre-setup/run