# configs/supervisord.conf
[supervisord]
nodaemon=true
user=abc
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
username=dummy
password=dummy

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
username=dummy
password=dummy

[program:cps]
command=python3 /app/calibre-web-automated/cps.py
directory=/app/calibre-web-automated
environment=CALIBRE_DBPATH="/config"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ingest-monitor]
command=/app/calibre-web-automated/scripts/service-wrappers/ingest-monitor.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:metadata-monitor]
command=/app/calibre-web-automated/scripts/service-wrappers/metadata-monitor.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:auto-zipper]
command=/app/calibre-web-automated/scripts/service-wrappers/auto-zipper.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:auto-library]
command=python3 /app/calibre-web-automated/scripts/auto_library.py
# Only restart if it fails? Original s6 script doesn't seem to loop, but if it's one-shot why is it a service?
# s6 run script: just runs python script. Using autorestart=unexpected with exitcodes=0
autostart=true
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
